# debbuild/debian/prerm
#!/bin/sh

systemctl disable jasper.service >/dev/null
systemctl stop jasper.service >/dev/null