# debbuild/debian/prerm
#!/bin/sh

set -e

deb-systemd-helper disable jasper.service >/dev/null
deb-systemd-invoke stop jasper.service >/dev/null
